# azubi_werkzeug/rootfs/etc/services.d/azubi-werkzeug/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Azubi Werkzeug Tracker
# Starts the Azubi Werkzeug Tracker application
# ==============================================================================

bashio::log.info "Starte Azubi Werkzeug Tracker..."

# Wechsel ins Arbeitsverzeichnis
cd /app || bashio::exit.nok "Kann nicht nach /app wechseln"

# Datenbank initialisieren
bashio::log.info "Initialisiere Datenbank..."
python3 -c "from app import setup_database; setup_database()" \
    || bashio::exit.nok "Datenbank-Initialisierung fehlgeschlagen"

# SSL Konfiguration lesen
SSL=$(bashio::config 'ssl')
CERTFILE=$(bashio::config 'certfile')
KEYFILE=$(bashio::config 'keyfile')

# Pruefe ob SSL aktiviert ist
if bashio::var.true "${SSL}"; then
    bashio::log.info "SSL ist aktiviert."
    bashio::log.info "Zertifikat: /ssl/${CERTFILE}"
    bashio::log.info "Key: /ssl/${KEYFILE}"
    
    # Server starten mit SSL
    bashio::log.info "Starte Gunicorn Server mit SSL auf Port 5000..."
    exec gunicorn \
        --workers 1 \
        --timeout 120 \
        --graceful-timeout 120 \
        --bind 0.0.0.0:5000 \
        --certfile "/ssl/${CERTFILE}" \
        --keyfile "/ssl/${KEYFILE}" \
        --access-logfile - \
        --error-logfile - \
        app:app
else
    bashio::log.info "SSL ist deaktiviert."
    
    # Server starten ohne SSL
    bashio::log.info "Starte Gunicorn Server auf Port 5000..."
    exec gunicorn \
        --workers 1 \
        --timeout 120 \
        --graceful-timeout 120 \
        --bind 0.0.0.0:5000 \
        --access-logfile - \
        --error-logfile - \
        app:app
fi
