# azubi_werkzeug/rootfs/etc/services.d/azubi-werkzeug/run
#!/bin/bash
# ==============================================================================
# Home Assistant Add-on: Azubi Werkzeug Tracker
# Starts the Azubi Werkzeug Tracker application
# ==============================================================================

echo "Starte Azubi Werkzeug Tracker..."

# Wechsel ins Arbeitsverzeichnis
cd /app || { echo "Kann nicht nach /app wechseln"; exit 1; }

# Datenbank initialisieren
echo "Initialisiere Datenbank..."
python3 -c "from app import setup_database; setup_database()" \
    || { echo "Datenbank-Initialisierung fehlgeschlagen"; exit 1; }

# Server starten
echo "Starte Gunicorn Server auf Port 5000..."
exec gunicorn \
    --workers 3 \
    --bind 0.0.0.0:5000 \
    --access-logfile - \
    --error-logfile - \
    app:app
