# azubi_werkzeug_beta/rootfs/etc/services.d/azubi-werkzeug/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Azubi Werkzeug Tracker
# Starts the Azubi Werkzeug Tracker application
# ==============================================================================

bashio::log.info "Starte Azubi Werkzeug Tracker..."

# Wechsel ins Arbeitsverzeichnis
cd /app || bashio::exit.nok "Kann nicht nach /app wechseln"

# Datenbank initialisieren (wird durch Gunicorn/app.py abgewickelt)
bashio::log.info "Initialisiere Datenbank..."

# SSL Konfiguration lesen
SSL=$(bashio::config 'ssl')
CERTFILE=$(bashio::config 'certfile')
KEYFILE=$(bashio::config 'keyfile')

# Pruefe ob SSL aktiviert ist
if bashio::var.true "${SSL}"; then
    bashio::log.info "SSL ist aktiviert."
    
    CERTFILE_PATH="/ssl/${CERTFILE}"
    KEYFILE_PATH="/ssl/${KEYFILE}"

    if [ ! -f "${CERTFILE_PATH}" ] || [ ! -f "${KEYFILE_PATH}" ]; then
        bashio::log.warning "Zertifikate (${CERTFILE_PATH} / ${KEYFILE_PATH}) nicht gefunden!"
        bashio::log.warning "Generiere temporäres (Self-Signed) Zertifikat für lokalen Zugriff..."
        openssl req -new -newkey rsa:4096 -days 3650 -nodes -x509 \
            -subj "/C=DE/ST=Berlin/L=Berlin/O=Local/CN=localhost" \
            -keyout "/data/adhoc_key.pem" \
            -out "/data/adhoc_cert.pem" 2>/dev/null
        CERTFILE_PATH="/data/adhoc_cert.pem"
        KEYFILE_PATH="/data/adhoc_key.pem"
        bashio::log.info "Temporäres Zertifikat erstellt: ${CERTFILE_PATH}"
    else
        bashio::log.info "Nutze vorhandene Zertifikate:"
        bashio::log.info "Zertifikat: ${CERTFILE_PATH}"
        bashio::log.info "Key: ${KEYFILE_PATH}"
    fi

    bashio::log.info "Konfiguriere NGINX als SSL-Reverse-Proxy auf Port 5000..."
    mkdir -p /etc/nginx/http.d
    cat > /etc/nginx/http.d/default.conf <<EOF
server {
    listen 5000 ssl;
    ssl_certificate ${CERTFILE_PATH};
    ssl_certificate_key ${KEYFILE_PATH};
    client_max_body_size 50m;
    location / {
        proxy_pass http://127.0.0.1:8099;
        proxy_set_header Host \$http_host;
        proxy_set_header X-Real-IP \$remote_addr;
        proxy_set_header X-Forwarded-For \$proxy_add_x_forwarded_for;
        proxy_set_header X-Forwarded-Proto \$scheme;
    }
}
EOF
else
    bashio::log.info "SSL ist deaktiviert."
    bashio::log.info "Konfiguriere NGINX als HTTP-Reverse-Proxy auf Port 5000..."
    mkdir -p /etc/nginx/http.d
    cat > /etc/nginx/http.d/default.conf <<EOF
server {
    listen 5000;
    client_max_body_size 50m;
    location / {
        proxy_pass http://127.0.0.1:8099;
        proxy_set_header Host \$http_host;
        proxy_set_header X-Real-IP \$remote_addr;
        proxy_set_header X-Forwarded-For \$proxy_add_x_forwarded_for;
        proxy_set_header X-Forwarded-Proto \$scheme;
    }
}
EOF
fi

bashio::log.info "Starte NGINX im Hintergrund..."
# suppress warning about running as root if any
nginx -g "daemon off;" 2>&1 &

bashio::log.info "Starte Gunicorn Server (Backend) auf Port 8099..."
exec gunicorn \
    --workers 1 \
    --threads 2 \
    --worker-class gthread \
    --timeout 120 \
    --graceful-timeout 120 \
    --keep-alive 5 \
    --max-requests 1000 \
    --max-requests-jitter 50 \
    --bind 0.0.0.0:8099 \
    --access-logfile - \
    --error-logfile - \
    app:app
